# templates/mariadb_maintenance/scripts/backup.cron
# Minute Hour DayOfMonth Month DayOfWeek Command

# Full backup daily at 00:00 AM
0 0 * * * /usr/local/bin/backup.sh full

# Incremental backups hourly on the hour
0 * * * * /usr/local/bin/backup.sh incremental

# Dump file hourly on the hour
# 0 * * * * /usr/local/bin/backup.sh dump