# templates/mariadb_maintenance/scripts/backup.cron
# Minute Hour DayOfMonth Month DayOfWeek Command

# Daily full backup at midnight
0 0 * * * /usr/local/bin/backup.sh full

# Incremental backup every hour on the hour
0 * * * * /usr/local/bin/backup.sh incremental

# Logical dump every hour (starts at minute 5 to avoid overlap)
#5 * * * * /usr/local/bin/backup.sh dump